# kodi-discdetect.service
[Service]
Type=oneshot
ExecStart=/usr/bin/kodi-detect-disc

